# orchestration/node-apps/hello-world/hello-relpath.service
[Unit]
Description=An app to say Hello World.

[Service]
Type=simple
User="${USER}"
WorkingDirectory=/home/"${USER}"/
ExecStartPre=npm install /home/"${USER}"/
ExecStart=node /home/"${USER}"/hello-world-app.js
# optional items below
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target